# single-node/provisioning/wazuh_manager/etc/rules/100206-trojan_snake.xml
<group name="snake,trojan">
  <!-- Windows Defender exclusion list -->
  <rule id="100206" level="12">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.commandLine" type="pcre2">powershell.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">Add-MpPreference\s+-ExclusionPath</field>
    <field name="win.eventdata.commandLine" type="pcre2">AppData\\\\Roaming\\\\.*exe</field>
    <description>Possible Snake keylogger activity: $(win.eventdata.parentImage) created a file in AppData Roaming folder and added it to the Windows Defender exclusion list.</description>
    <mitre>
      <id>T1087</id>
      <id>T1059.003</id>
    </mitre>
  </rule>

  <!-- Scheduled task created -->
  <rule id="100207" level="12">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">schtasks.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">\/Create\s+\/TN</field>
    <field name="win.eventdata.commandLine" type="pcre2">\/XML</field>
    <field name="win.eventdata.commandLine" type="pcre2">Updates</field>
    <field name="win.eventdata.commandLine" type="pcre2">AppData\\\\Local\\\\Temp\\\\[^\\]+\.tmp</field>
    <description>Possible Snake keylogger activity: Scheduled task created.</description>
    <mitre>
      <id>T1055.001</id>
    </mitre>
  </rule>

  <!-- Registry key created -->
  <rule id="100208" level="15">
    <if_sid>61614</if_sid>
    <field name="win.eventdata.EventType" type="pcre2">CreateKey$</field>
    <field name="win.eventdata.Image" type="pcre2">svchost.exe</field>
    <field name="win.eventdata.TargetObject" type="pcre2">WeENKtk</field>
    <description>Possible Snake keylogger activity: Suspicious registry key created.</description>
    <mitre>
      <id>T1204.002</id>
      <id>T1547.001</id>
    </mitre>
  </rule>

  <!-- Registry subkeys created -->
  <rule id="100209" level="15">
    <if_sid>61615</if_sid>
    <field name="win.eventdata.eventType" type="pcre2">SetValue$</field>
    <field name="win.eventdata.targetObject" type="pcre2">HKLM\\\\SOFTWARE\\\\Microsoft\\\\Windows NT\\\\CurrentVersion\\\\Schedule\\\\TaskCache\\\\Tree\\\\Updates\\\\[^\\]+\\\\(SD|Id|Index)</field>
    <description>Possible Snake keylogger activity: Suspicious registry subkey creation associated with Snake keylogger.</description>
    <mitre>
      <id>T1543</id>
    </mitre>
  </rule>

  <!-- File deletion -->
  <rule id="100210" level="15">
    <if_sid>61603</if_sid>
    <field name="win.eventdata.image" type="pcre2">(?i)choice\.exe</field>
    <field name="win.eventdata.commandLine" type="pcre2">choice\s+\/C\s+Y\s+\/N\s+\/D\s+Y\s+\/T\s+3</field>
    <description>Possible Snake keylogger activity: Snake keylogger deleted.</description>
    <mitre>
      <id>T1055.001</id>
    </mitre>
  </rule>
</group>